# files/refresh-termstats.timer
[Unit]
Description=Refresh terminals statistics every six hours

[Timer]
OnCalendar=0/6:00

[Install]
WantedBy=multi-user.target
